# config/99-battery_notification.rules
SUBSYSTEM=="power_supply", ACTION=="change", \
RUN+="/bin/su -c '/usr/local/bin/power_handler.sh' - $SUDO_USER"